# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command = bandersnatch mirror
startsecs = 5
stopwaitsecs = 3600
stopasgroup = false
killasgroup = true
autorestart=true
